# upgrade/update/service/vfenv.service
[Unit]
Description=vfenv

[Service]
ExecStart=/usr/sbin/vfenv.sh

[Install]
WantedBy=multi-user.target
